# resources/autologin/etc/systemd/system/getty@tty1.service.d/override.conf.template
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin ${MINER_USER} %I $TERM
